# alarm.timer
[Unit]
Description=Daily wake up alarm

[Timer]
OnCalendar=Mon-Fri *-*-* 6:45:00

[Install]
WantedBy=timers.target
