# alarm-server.service
[Unit]
Description=Alarm Clock Server

[Service]
ExecStart=alarm-server

[Install]
WantedBy=multi-user.target
