# alarm.service
[Unit]
Description=Alert Trigger

[Service]
ExecStart=alarm start

[Install]
WantedBy=multi-user.target
