# oncall-monitor.service
[Unit]
Description=On-Call Alert Monitor

[Service]
ExecStart=oncall-monitor

[Install]
WantedBy=multi-user.target
